# supervisor-app.conf
[program:app-uwsgi]
command = uwsgi --ini /appdata/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

[supervisord]
loglevel=debug